# agents/agent_one/agent_one.service
[Unit]
Description=invinoveritas Agent One autonomous buyer agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/invinoveritas/agents/agent_one
Environment=INVINO_API_BASE=https://api.babyblueviper.com
Environment=AGENT_ONE_DATA_DIR=/var/lib/agent_one
Environment=AGENT_ONE_LOG_FILE=/var/log/agent_one.log
Environment=AGENT_ONE_MIN_BALANCE_SATS=5000
Environment=AGENT_ONE_FIRST_BUY_MAX_SATS=1500
Environment=AGENT_ONE_DAILY_BUY_LIMIT=2
Environment=AGENT_ONE_DAILY_SPEND_CAP_SATS=5000
ExecStart=/root/invinoveritas/venv/bin/python /root/invinoveritas/agents/agent_one/agent_one.py
Restart=always
RestartSec=30
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
